# roles/php/templates/etc/systemd/system/php8.3-fpm.service.d/symbiotic.conf
# {{ ansible_managed }}
#
[Service]
# Always restart, because after an OOM kill, systemd does not always restart
Restart=always
