# roles/mariadb/templates/etc/systemd/system/mariadb.service.d/symbiotic.conf
# {{ ansible_managed }}
#
[Service]
# Avoid open file limits. See coopsymbiotic/ops#633
LimitNOFILE=196608

# Always restart, because after an OOM kill, systemd does not always restart
Restart=always
# Default is 100ms
RestartSec=20s

# Kill other processes before killing mysql (which might free whatever is clogging mysql)
OOMScoreAdjust=-600
