# roles/mariadb/templates/etc/systemd/system/mariadb.service.d/symbiotic.conf
# {{ ansible_managed }}
#
[Service]
# Avoid open file limits.
LimitNOFILE=65536

# Always restart, because after an OOM kill, systemd does not always restart
Restart=always
