# roles/mariadb/templates/etc/systemd/system/mariadb.service.d/symbiotic.conf
# {{ ansible_managed }}
#
# Avoid open file limits.
[Service]
LimitNOFILE=65536
