# roles/civicrm-sites/templates/etc/systemd/system/d7civicron.timer
# {{ ansible_managed }}

[Unit]
Description="D7-CiviCRM Cron for {{ inventory_hostname }}"
Requires=d7civicron_{{ site_shortname }}.service

[Timer]
Unit=d7civicron_{{ site_shortname }}.service
OnUnitInactiveSec=15min
OnBootSec=10min
RandomizedDelaySec=30sec

[Install]
WantedBy=timers.target
